# deploy/jy-trade-api.service.example
[Unit]
Description=JY Trade API
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/jy-trade/current
EnvironmentFile=/opt/jy-trade/current/.env
Environment=TZ=Asia/Shanghai
ExecStart=/usr/bin/npm run start:api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
